# launchd/com.silopolis.deploy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  SILOPOLIS — Auto-commit to GitHub (NOT direct Vercel deploy)
  =============================================================
  Runs auto_deploy.sh every 30 minutes.
  • If dashboard/, core/, or api/ have uncommitted changes: git commit + push
  • Vercel GitHub integration handles production builds automatically on push
  • This script NEVER calls `vercel` directly

  StartInterval: 1800 → check every 30 minutes
  RunAtLoad: false → do NOT deploy immediately on load

  Logs:
    tail -f /tmp/silopolis-deploy.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.silopolis.deploy</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd "$(dirname "$0")/.." &amp;&amp; bash scripts/auto_deploy.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/funmen/silopolis</string>

    <!-- Check every 30 minutes — Vercel deploys on git push, not on every run -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <key>StandardOutPath</key>
    <string>/tmp/silopolis-deploy.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/silopolis-deploy.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/funmen</string>
    </dict>
</dict>
</plist>
