# launchd/com.silopolis.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.silopolis.tunnel</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/funmen/silopolis/scripts/start_tunnel.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/silopolis-tunnel.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/silopolis-tunnel.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
