# launchd/com.silopolis.heartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  SILOPOLIS macOS LaunchAgent — 12 heartbeats/day (every 2 hours)
  Active mode for 5-day hackathon sprint.

  Install:
    cp launchd/com.silopolis.heartbeat.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.silopolis.heartbeat.plist

  View logs:
    tail -f /tmp/silopolis-heartbeat.log

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.silopolis.heartbeat.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.silopolis.heartbeat</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/funmen/silopolis/core/heartbeat.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/funmen/silopolis</string>

    <!-- Fire every 2 hours: 00:00, 02:00, 04:00, 06:00, 08:00, 10:00,
                             12:00, 14:00, 16:00, 18:00, 20:00, 22:00 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>0</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>2</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>4</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>20</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>22</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/tmp/silopolis-heartbeat.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/silopolis-heartbeat.err</string>

    <key>KeepAlive</key>
    <false/>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
